# bernard/harness/dot-claude/scripts/com.lashawsalta.gbrain-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.lashawsalta.gbrain-sync</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>bash /Users/lashawsalta/life-wiki/.claude/scripts/gbrain-job.sh sync</string>
  </array>
  <!-- gbrain-job.sh writes the newest-first dated run log to gbrain-sync.log
       and filters the [extract.*] progress-ticker spam. The paths below only
       catch stray wrapper output. -->

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>3</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/lashawsalta/life-wiki/.claude/logs/gbrain-sync-launchd.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/lashawsalta/life-wiki/.claude/logs/gbrain-sync-launchd.err</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/Users/lashawsalta/.bun/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>
</dict>
</plist>
