# bernard/harness/dot-claude/scripts/com.life-wiki.synthesis.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.life-wiki.synthesis</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>/Users/lashawsalta/life-wiki/.claude/scripts/daily-synthesis.sh</string>
  </array>

  <!-- Daily 06:00 — after gbrain dream (03:15) and the leak scan (02:30).
       Cheap on quiet days: the script's snapshot gate exits before invoking
       Claude when nothing changed since the last completed pass. Runs on next
       wake if the Mac was asleep.
       On-demand: bash .claude/scripts/daily-synthesis.sh -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>6</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/lashawsalta/life-wiki/.claude/logs/daily-synthesis-launchd.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/lashawsalta/life-wiki/.claude/logs/daily-synthesis.err</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>
</dict>
</plist>
